# script/cron/checkin.cron
30 8 * * 1-5 /bin/bash /var/www/guodaxia/script/cron/checkin.sh
